# crontab
0 16 * * * python /media-cleaner/media_cleaner.py >> /media-cleaner/cleaner.log